# target/linux/imx/imx8/base-files/etc/uci-defaults/13_firewall_4g
. /lib/functions.sh
# add 4G to firewall WAN zone

NET='wwan0'

uci get firewall.@zone[1].network | grep -q "$NET" || {
	uci -q batch <<-EOF >/dev/null
		add_list firewall.@zone[1].network="$NET"
		commit firewall
	EOF
}

exit 0
